# huion-note-x10.service
[Unit]
Description=Huion Note X10 BLE tablet driver
After=bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.local/share/huion-note-x10/huion_ble_driver.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
